# data/io.github.XZram.service
[Unit]
Description=XZram swap management D-Bus daemon
After=dbus.service

[Service]
Type=dbus
BusName=io.github.XZram1
ExecStart=/usr/libexec/xzramd
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
# Apply/restore run xzram-helper in a transient unit (see xzramd privileged.rs).
ReadWritePaths=/var/lib/xzram
SystemCallFilter=@system-service

[Install]
WantedBy=multi-user.target
